# internal/render/testdata/golden/clusters/d7a/nats-server.service
# AUTOGENERATED by vikasa-infra/cmd/gen
[Unit]
Description=Vikasa NATS server (JetStream + leaf)
Documentation=https://github.com/Vikasa2M/vikasa-infra
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nats
Group=nats
ExecStart=/usr/local/bin/nats-server -c /etc/vikasa/nats.conf
Restart=on-failure
RestartSec=5
StateDirectory=nats
StateDirectoryMode=0750
ReadWritePaths=/var/lib/nats
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
